# docs/supervisor/wolf-reverb.conf
; Deploy instructions:
; 1. Copy this file to /etc/supervisor/conf.d/wolf-reverb.conf on the VPS
; 2. Run: supervisorctl reread
; 3. Run: supervisorctl update
; 4. Run: supervisorctl start wolf-reverb
; 5. Check status: supervisorctl status wolf-reverb
;
; Adjust 'command' path and 'user' to match your VPS setup.

[program:wolf-reverb]
command=php /var/www/wolf/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/wolf
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/wolf/storage/logs/reverb.log
