# docs/supervisor/wolf-queue.conf
; Deploy instructions:
; 1. Copy this file to /etc/supervisor/conf.d/wolf-queue.conf on the VPS
; 2. Run: supervisorctl reread
; 3. Run: supervisorctl update
; 4. Run: supervisorctl start wolf-queue
; 5. Check status: supervisorctl status wolf-queue
;
; After deploying new code, restart the worker so it picks up changes:
;   supervisorctl restart wolf-queue
;
; Adjust 'command' path and 'user' to match your VPS setup.

[program:wolf-queue]
command=php /var/www/wolf/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/wolf
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/wolf/storage/logs/queue.log
stopwaitsecs=3600
